# AI/AI-Coding/LiteLLM-SearXNG-AICoding-Gateway-Single-ECS/assets/config/litellm.service.example
[Unit]
Description=LiteLLM Proxy for Huawei Cloud MaaS
After=network-online.target redis-server.service postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=litellm
Group=litellm
WorkingDirectory=/opt/litellm
EnvironmentFile=/etc/litellm/litellm.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/litellm-venv/bin/litellm --config /etc/litellm/config.yaml --host 0.0.0.0 --port 4000 --use_prisma_db_push
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
